# openwrt-23.05/files/etc/hotplug.d/iface/06-ethenableflowcontrol
#!/bin/sh
logger "enabling-eth-flow-control-from-hotplug"
ethtool -A eth0 autoneg on rx on tx on
ethtool -A eth1 autoneg on rx on tx on
